# src/kubic-archive-DVD.service.in
[Unit]
Description=Backup the openSUSE Kubic Installer DVD
Wants=network.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/kubic-archive/mirror-openSUSE-Kubic-DVD
IOSchedulingClass=best-effort
IOSchedulingPriority=7
