# src/kubic-archive-container.service.in
[Unit]
Description=Backup the openSUSE registry into a local registry
Wants=network.target
After=network-online.target

[Service]
Type=oneshot
Environment="HOME=/var/lib/empty"
ExecStart=@libexecdir@/kubic-archive/mirror-openSUSE-registry
